# 1wirevz
### BEGIN INIT INFO
# Provides: 1wirevz 
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 1
# Short-Description: DS2482 I²C 1-Wire® Master to Volkszaehler RaspberryPI deamon
### END INIT INFO

#!/bin/bash
 
DAEMON="/usr/sbin/1wirevz"
PIDFILE="/tmp/1wirevz.pid" 

test -f $DAEMON || exit 0
 
case "$1" in

   start)
	test ! -f $PIDFILE || { echo "Deamon already running!"; exit 0; }
        $DAEMON
        ;;

    stop)
        killall -q 1wirevz
	rm -f /tmp/1wirevz.pid
        ;;

    restart)
        $0 stop
        $0 start
        ;;

    *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
        ;;
esac
